# transport_display_config.service
[Unit]
Description=Bus departure display config web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs unprivileged; restarting the display service is allowed via a narrow
# sudoers drop-in installed by setup_pi.sh (see /etc/sudoers.d/).
User=pierre
# WorkingDirectory and ExecStart are rewritten by setup_pi.sh to match where the
# repo is cloned and its web venv (webenv/) that holds FastAPI/uvicorn.
WorkingDirectory=/home/pierre/transport_display
ExecStart=/home/pierre/transport_display/webenv/bin/python3 -m server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
