# transport_display.service
[Unit]
Description=Bus departure LED matrix display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The matrix library needs root for GPIO/PWM; it drops privileges after init
# (RGBMatrixOptions.drop_privileges = True).
User=root
# WorkingDirectory and ExecStart are rewritten by setup_pi.sh to match where the
# repo is cloned and its venv (env/) that holds the rgbmatrix bindings.
WorkingDirectory=/home/pierre/transport_display
ExecStart=/home/pierre/transport_display/env/bin/python3 -m src
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
