# transport_display_config.service
[Unit]
Description=Bus departure display config web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs unprivileged; restarting the display service is allowed via a narrow
# sudoers drop-in installed by setup_pi.sh (see /etc/sudoers.d/).
# User, WorkingDirectory, and ExecStart are placeholders rewritten by
# setup_pi.sh to match the chosen SERVICE_USER, where the repo is cloned, and
# its web venv (webenv/) that holds FastAPI/uvicorn.
User=pi
WorkingDirectory=/home/pi/transport_display
ExecStart=/home/pi/transport_display/webenv/bin/python3 -m server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
